# chroot/service/weave/run
#!/bin/sh
exec 2>&1
set -e
docker history zettio/weave
weave launch
exec docker wait weave
