# chroot/service/prepare_pciids/run
#!/bin/sh
if grep -q 'heddle_prepare=1' /proc/cmdline; then
  exec > /dev/ttyS0 2>&1
  update-pciids
else
  exec 2>&1
fi
exec sv down prepare_pciids
