# chroot/service/dhcpcd/run
#!/bin/sh
exec 2>&1
rm -f /var/run/dhcpcd.pid
exec dhcpcd -B -Z 'dummy*,docker*,weave,veth*'
