# xroot/service/ntpd/run
#!/bin/sh
exec 2>&1
peers=(0.heddle.pool.ntp.org 1.heddle.pool.ntp.org 2.heddle.pool.ntp.org 3.heddle.pool.ntp.org)
exec ntpd -n ${peers[@]/#/-p }
