# chroot/service/agetty-tty2/run
#!/bin/sh
exec setsid agetty --nohostname -o '-p -- \u' tty2
