# chroot/service/agetty-tty1/run
#!/bin/sh
exec setsid agetty -o '-p -- \u' tty1
