# chroot/service/docker/run
#!/bin/sh
exec 2>&1
mkdir -p /var/lib
ln -sf /extra/docker /var/lib
rm -f /var/run/docker.pid
exec docker -d
